# scheduled/launchd/io.lorecraft.mogging.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>io.lorecraft.mogging.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-c</string>
        <string>cd "$VAULT" &amp;&amp; claude -p 'execute agents/morning.md in audit mode' --headless --audit</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>VAULT_ROOT</key>
        <string>$VAULT</string>
        <key>REPO_ROOT</key>
        <string>$REPO_ROOT</string>
        <key>PATH</key>
        <string>$HOME/.local/bin:$HOME/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>$HOME/Library/Logs/mogging-morning.log</string>
    <key>StandardErrorPath</key>
    <string>$HOME/Library/Logs/mogging-morning.log</string>

    <key>RunAtLoad</key>
    <false/>
    <key>KeepAlive</key>
    <false/>
    <key>AbandonProcessGroup</key>
    <true/>
    <key>ProcessType</key>
    <string>Background</string>
</dict>
</plist>
